# blocklist-update.service
[Unit]
Description=Updates dns block lists for dnsmasq

[Service]
Type=simple
ExecStart=/usr/local/bin/blocklist-update
